# etc/bot.service
[Unit]
Description=Hexapod

[Service]
ExecStart=/root/bot
User=root
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
