# etc/bot.service
[Unit]
Description=Hexapod

[Service]
ExecStart=/home/pi/bot
User=root
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
