# etc/sixad.service
[Unit]
Description=Sixaxis Bluetooth Interface

[Service]
ExecStart=/usr/bin/sixad
User=root
Restart=always

[Install]
WantedBy=multi-user.target
